# fortune@.service
[Unit]
Description=QoTD

[Service]
ExecStart=-/usr/bin/fortune
StandardOutput=socket
DynamicUser=yes
ProtectHome=true
PrivateUsers=true